# dist/nvr.service
[Unit]
Description=Hybrid NVR Service
After=network.target

[Service]
Type=simple
User=nvr
Group=nvr
WorkingDirectory=/opt/nvr/bin
ExecStart=/opt/nvr/bin/nvr_service
Restart=always
RestartSec=5
# Allow unrestricted access to /dev/shm for the shared memory ring buffer
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target